# linux-control-client/linux-control-client.install
post_install() {
    cat <<EOF
Make sure you are in the power group so you can shutdown, restart, etc.
and then edit /etc/linux-control-client.conf. Finally, run:
    pip install --user plocate python-libxdo    
    sudo systemctl restart polkit
    sudo systemctl enable linux-control-client
    sudo systemctl restart linux-control-client
EOF
}
